# src/JenkinsConfig/isas.service
# 信息智能分析系统 systemd 服务单元 (Linux)
#
# 端口由 config/app.json 的 server.port 决定（run.sh 读取），无需在本单元硬编码。
# 如需覆盖，取消下行注释并改端口：
#   Environment=ISAS_SERVER_PORT=29000
#
# 安装方式（部署路径以 /opt/InformationSmartAnalysisSystem 为例，按实际调整）：
#   sudo cp src/JenkinsConfig/isas.service /etc/systemd/system/isas.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now isas
# 常用命令：
#   sudo systemctl start isas | stop isas | restart isas | status isas
#   sudo journalctl -u isas -f
# 注意：部署目录属 jenkins 用户，故以 jenkins 运行；如改用其它用户请同步 chown。

[Unit]
Description=InformationSmartAnalysisSystem (信息智能分析系统)
After=network.target

[Service]
Type=simple
User=jenkins
WorkingDirectory=/opt/InformationSmartAnalysisSystem
ExecStart=/opt/InformationSmartAnalysisSystem/run.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
